# dotfiles/.config/service/stw-service/run
#!/bin/sh

				exec stw -p 1 -B 10 -y 10 psvstat -s /var/service ~/.xservice/
			
