# dotfiles/.config/service/dwm/run
#!/bin/sh

				cd $HOME
				exec dwm
			
