# dotfiles/.config/service/stacki3/run
#!/bin/bash

cd ~/stacki3
source venv/bin/activate

exec stacki3 45
